# scripts/etc/init.d/na62-startup
#!/bin/bash
# chkconfig: 5 99 90
/performance/load-pf_ring

#ssh -f na62farmui1 -L 2000:github.com:22 -N
#ssh -f na62farmui1 -L 2001:etap-git.physik.uni-mainz.de:22 -N
